# scheduler/com.sermon-study.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sermon-study.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/johau/Projects/sermon-study/node_modules/.bin/tsx</string>
    <string>/Users/johau/Projects/sermon-study/scripts/daily-study.ts</string>
  </array>

  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Hour</key>
      <integer>21</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Hour</key>
      <integer>22</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/johau/Projects/sermon-study</string>

  <key>StandardOutPath</key>
  <string>/Users/johau/Projects/sermon-study/data/daily.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/johau/Projects/sermon-study/data/daily.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
  </dict>
</dict>
</plist>
